# debian-11-PrisonPC/xfce/popcon.timer
[Unit]
Description=popularity contest

[Timer]
# Every ten seconds.
OnCalendar=*:*:0/10

[Install]
WantedBy=timers.target
