# debian-11-PrisonPC/xfce/thunar-bookmark-srv-share.service
[Service]
Type=oneshot
ExecStart=bootstrap2020-thunar-bookmark-srv-share.py
[Install]
WantedBy=default.target
