# debian-11-main/cgit.socket
[Unit]
PartOf=cgit.service

[Socket]
ListenStream=/run/cgit.socket

# Allow nginx (fastcgi-pass) to connect to this AF_UNIX socket.
SocketUser=www-data
SocketGroup=www-data
SocketMode=0660

[Install]
WantedBy=sockets.target
