# debian-11-datasafe3/dyndns.timer
[Timer]
OnBootSec=1min
OnUnitActiveSec=6h

[Install]
WantedBy=timers.target
