# debian-12-datasafe3.files/rsnapshot.service
[Service]
Type=oneshot
ExecStart=/usr/local/bin/datasafe-rsnapshot-wrapper
