# debian-12-datasafe3.files/storage-check.service
[Service]
Type=oneshot
ExecStart=/usr/local/bin/storage-check
