# debian-11-PrisonPC-tvserver/systemd/record.service
[Service]
Type=oneshot
ExecStart=tvserver-record
