# debian-11-PrisonPC-tvserver/systemd/update-config.service
[Service]
Type=oneshot
ExecStart=tvserver-update-config
