# debian-12-PrisonPC-tvserver.files/systemd/check_tv.timer
[Timer]
OnCalendar=*:0/2
AccuracySec=1m
[Install]
WantedBy=timers.target
