# debian-12-PrisonPC-tvserver.files/systemd/epg-scan@.service
[Service]
Type=oneshot
ExecStart=tvserver-epg-scan --adapter %I
