# debian-12-PrisonPC-tvserver.files/systemd/import-media.service
[Service]
Type=oneshot
# This was "iptv import-media"
ExecStart=tvserver-import-media
