# debian-12-PrisonPC-tvserver.files/systemd/multicat@.service
[Service]
Restart=always
RestartSec=30s
StartLimitBurst=0
ExecStart=multicat -U @%I:1234 /dev/null
