# debian-12-PrisonPC-tvserver.files/systemd/record.timer
[Timer]
# NOTE: "minutely" means "every minute", i.e. what cron calls "* *  * * *".
OnCalendar=minutely
[Install]
WantedBy=timers.target
