# debian-11-PrisonPC/snitches/usb-snitch.service
[Unit]
DefaultDependencies=no
Before=udev.service

[Service]
Type=notify
ExecStart=/usr/local/sbin/usb-snitchd

FailureAction=reboot
