# debian-11-PrisonPC/pki-nssdb/nssdb-update.service
[Service]
Type=oneshot
ExecStart=bootstrap2020-nssdb-install.py
[Install]
WantedBy=default.target
