# debian-11-PrisonPC/srv-share.mount
[Mount]
What=nfs:/srv/share
Where=/srv/share
Type=nfs
Options=nodev,noexec,nosuid,ro
DirectoryMode=0

[Install]
WantedBy=remote-fs.target
